# init/trackermon.service
[Unit]
Description=Trackermon test

[Service]
Type=simple
User=root
ExecStart=/home/ferllini13/Desktop/Operativos/short1/TrackerMon-Daemon/bin/trackermon

[Install]
WantedBy=multi-user.target
