# init/trackermon.service
[Unit]
Description=Trackermon test

[Service]
Type=simple
User=root
ExecStart=/etc/trackermon/trackermon

[Install]
WantedBy=multi-user.target